# roles/webapp/templates/foreman.service.j2
# {{ansible_managed}}

[Unit]
Description=Foreman app service

[Service]
ExecStart=/bin/bash -lc "cd ~/app/current && BIND={{foreman_bind}} foreman start"

[Install]
WantedBy=multi-user.target
